# rootfs/etc/s6-overlay/s6-rc.d/configsetup/up
/etc/s6-overlay/scripts/configchangescript
