# rootfs/etc/s6-overlay/s6-rc.d/status/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: WireGuard
# Shows current WireGuard status
# ==============================================================================
sleep 60
bashio::log.info "---------Current WireGuard status---------"

data=$(exec wg show)

public_key=$(echo "$data" | grep "public key:" | awk '{print $NF}')
latest_handshake=$(echo "$data" | grep "latest handshake:" | awk '{$1=$2=""; print $0}')
transfer=$(echo "$data" | grep "transfer:" | awk '{$1=""; print $0}')

bashio::log.info "Last handshake:$latest_handshake"
bashio::log.info "Data transfer:$transfer"
bashio::log.info "Public key:$public_key"
bashio::log.info "-------------End-------------"