# rootfs/etc/s6-overlay/s6-rc.d/backup/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: m-connect
# Runs backup tasks
# ==============================================================================
#https://linuxize.com/post/how-to-create-users-in-linux-using-the-useradd-command/?utm_content=cmp-true

# Get the supervisor token? Maybe auth with username and password?

#curl -sSL -H "Authorization: Bearer <<BEARERTOKEN>>" http://supervisor/backups/info

#curl -sSL -H "Authorization: Bearer <<BEARERTOKEN>>" -X POST http://supervisor/backups/new/full

#Result = {"result": "ok", "data": {"slug": "3d8be061"}}
#Sjekk for ok og sett opp synch, hvis ikke send feilmelding i loggen.

#Remote backup_location="gapit@10.0.0.23:/home/gapit/production_backup/"
#Have slugname as variable
curl -sSL -H "Authorization: Bearer <<BEARERTOKEN>>" http://supervisor/backups/a52292b1/download > /data/m-connect/backup/ha_backup_acrobatic-sailfish
rsync /data/m-connect/backup/ha_backup_acrobatic-sailfish acrobatic-sailfish@host-ip:/m-connect/acrobatic-sailfish
