# rootfs/etc/s6-overlay/s6-rc.d/setupbackup/up
/etc/s6-overlay/scripts/startupbackupscript
